# .payload/etc/cron.d/update-dnsmasq.sh
1    1    *    *    *    root    /config/scripts/update-dnsmasq -f /config/user-data/blacklist.cfg >/dev/null 2>&1